# nifty_pipeline.service
[Unit]
Description=Nifty Pipeline Trading Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shubham
WorkingDirectory=/home/shubham/nse_scheduler
ExecStart=/home/shubham/nse_scheduler/venv/bin/python /home/shubham/nse_scheduler/app.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
